# linux/systemd/i3.target
[Unit]
Description=i3 target
PartOf=default.target
